# k8s-snap/etc-systemd-system/snap.k8s.k8sd.service.d/override.conf
[Unit]
Requires=clean-up-k8s.service
After=clean-up-k8s.service
